# docker/hermes-agent.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/hermes-supervisord.pid
childlogdir=/tmp

[program:gateway]
command=%(ENV_HERMES_HOME)s/venv/bin/hermes gateway run
directory=/AMP
autostart=true
autorestart=true
startsecs=3
startretries=3
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=40

[program:dashboard]
command=%(ENV_HERMES_HOME)s/venv/bin/hermes dashboard --host 0.0.0.0 --port %(ENV_HERMES_DASHBOARD_PORT)s --no-open
directory=/AMP
autostart=true
autorestart=true
startsecs=3
startretries=3
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=40
